# etc/profile.d/condor.sh
#!/bin/bash
#
# Profile setter for condor utils
#
# Author: Duncan Macleod <duncan.macleod@ligo.org>
#

# aliases
alias cq="condor_q $USER"
alias crm="condor_rm"
alias cqnb="condor_q $USER -nobatch"
alias csub="condor_submit"
