# etc/profile.d/conda.sh
#!/bin/sh
#
# profile setter for conda

CONDA_PATH=${HOME}/anaconda2
alias conda="${HOME}/anaconda2/bin/conda"
