# etc/profile.d/core.sh
#
# Profile setter for core utils
#
# Author: Duncan Macleod <duncan.macleod@ligo.org>
#

# aliases
alias rm='rm -i'
alias view="vim -R"

# environment
export CLICOLOR=1
export EDITOR=vim
