# etc/profile.d/git.sh
#
# Profile setter for git
#
# Author: Duncan Macleod <macleoddm@cardiff.ac.uk>
#

if [ -f /usr/share/bash-completion/completions/git ]; then
  . /usr/share/bash-completion/completions/git
fi
