# etc/profile.d/grep.sh
#
# Profile setter for python
#
# Author: Duncan Macleod <duncan.macleod@ligo.org>
#

alias grep="grep --color=auto"
