# etc/profile.d/python.sh
#
# Profile setter for python
#
# Author: Duncan Macleod <duncan.macleod@ligo.org>
#

# always show DeprecationWarnings
export PYTHONWARNINGS="always::DeprecationWarning"
